# REST/scrapers_deployment/crontab
# run everyday at midnight : 0 0 * * * <command>
0 */2 * * * python3 /usr/src/pingodoce.py >> /usr/src/cron.log ; python3 /usr/src/auchan.py >> /usr/src/cron.log

